# imagenes/bind-dns-parcial/con composer/keepalive/files/init.keepalived.conf
description "Ip flotante para dns de parcial 1 protocolos de comunicacion"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec /usr/local/sbin/keepalived --dont-fork
